# src/resources/install-helpers/systemd/beanie-password-setup.service
[Unit]
Description=BeanieDeploy Password Setup GUI
After=graphical.target
Wants=graphical.target

[Service]
Type=simple
ExecStart=/usr/local/bin/beanie-password-setup
User=root
StandardOutput=journal
StandardError=journal
Restart=on-failure
RestartSec=5s
RemainAfterExit=yes

# Exit if marker file exists (password already set)
ExecStartPre=/bin/bash -c 'test ! -f /etc/beanie_password_set'

[Install]
WantedBy=graphical.target
