# src/resources/install-helpers/partition-resize/wingone-partition-resize.service
[Unit]
Description=Post-Install Partition Resize
# Run after local filesystems are mounted (root btrfs is up, online resize is safe),
# but before sysinit.target so nothing else is competing for the disk yet.
DefaultDependencies=no
After=local-fs.target
Before=sysinit.target

[Service]
Type=oneshot
ExecStart=/usr/bin/python3 /usr/local/bin/wingone-partition-resize.py
# The script disables itself on success.
ExecStartPost=-/usr/bin/systemctl disable wingone-partition-resize.service
ExecStartPost=-/usr/bin/rm -f /usr/local/bin/wingone-partition-resize.py 
ExecStartPost=-/usr/bin/rm -f /var/usrlocal/bin/wingone-partition-resize.py
ExecStartPost=-/usr/bin/rm -f /etc/systemd/system/wingone-partition-resize.service

[Install]
WantedBy=sysinit.target
